# backup/Brgy.Tanod-S.O.S-main/scripts/brgy-tanod.service
[Unit]
Description=Brgy Tanod SOS System
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/youruser/Brgy.Tanod-S.O.S
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
